# docker/etc/crontabs/root
* * * * * ( restart_services )
* * * * * ( sleep 10; restart_services )
* * * * * ( sleep 20; restart_services )
* * * * * ( sleep 30; restart_services )
* * * * * ( sleep 40; restart_services )
* * * * * ( sleep 50; restart_services )