# docker/etc/crontabs/root
* * * * * ( /app/restart.services.sh )
* * * * * ( sleep 10;/app/restart.services.sh )
* * * * * ( sleep 20;/app/restart.services.sh )
* * * * * ( sleep 30;/app/restart.services.sh )
* * * * * ( sleep 40;/app/restart.services.sh )
* * * * * ( sleep 50;/app/restart.services.sh )